# etc/supervisor/eventhandler.conf
[program:eventhandlers]
user=root
directory=/.tutorial
command=python3 eventhandlers.py
