# etc/supervisor/eventhandler.conf
[program:eventhandler]
user=root
directory=/home/user/tutorial
command=python3 app.py
