# etc/supervisor/eventhandler.conf
[program:eventhandler]
user=root
directory=/.tutorial
command=python3 app.py
