# etc/supervisor/terminado_ssh.conf
[program:sshd]
command=/usr/sbin/sshd -Df /etc/ssh/sshd_config_tfw